# crontab
1 * * * * sudo sh /etc/start/doi.sh;
10 3 * * 2 sudo sh /etc/start/yo.sh;
@reboot sudo iptables -t mangle -A PREROUTING -p IGMP -j DROP;
